# prebuilt/common/etc/cron/cron.daily/01clear_cache
#!/system/bin/sh

. /system/etc/cron/cron.conf

if $daily_clear_data_cache
then
	busybox find /data/data -type d -iname "*cache*" -maxdepth 2 -mindepth 2 -exec busybox rm -rf {} ';'
fi