# prebuilt/common/etc/cron/cron.weekly/00drop_caches
#!/system/bin/sh

. /system/etc/cron/cron.conf

if $weekly_drop_caches
then
	sync
	echo 3 > /proc/sys/vm/drop_caches
fi